# etc/systemd/user/gabi-eyecare.service
[Unit]
Description= For you, I love you so much ! <3

[Service]
Type=simple
ExecStart=/usr/local/bin/gabi-eyecare
